# 99-k480.rules
KERNEL=="hidrawXX", SUBSYSTEM=="hidraw", ACTION=="add", SUBSYSTEMS=="bluetooth", RUN+="/opt/bin/k480.sh %p"